# s6/.s6/sv/dbus-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
multisubstitute {
  importas ADDR DBUS_SESSION_BUS_ADDRESS
  importas HOME HOME
}

dbus-daemon --nofork --session --print-pid 4 --address=${ADDR} 
