# s6/.local/share/s6/sv/xdg-desktop-portal-log/run
#!/usr/bin/execlineb -P

s6-logscript
xdg-desktop-portal
