# s6/.s6/sv/feh/up
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas HOME HOME

sh ${HOME}/.s6/sv/feh/shell_up
