# s6/.s6/sv/xorg-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
multisubstitute {
    importas home HOME
    importas xauthority XAUTHORITY
}

backtick -E tty { pipeline { tty } tr -d '/dev/tty/' }

foreground { touch $xauthority }
foreground { stty sane }
foreground { xauth remove :$tty }

backtick -E cookie { pipeline { od -An -N16 -tx /dev/urandom } tr -d " " }
foreground { xauth add :$tty MIT-MAGIC-COOKIE-1 $cookie }
Xorg :$tty -keeptty vt$tty -noreset -auth $xauthority -displayfd 4
