# s6/.s6/sv/xorg-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas XAUTHORITY XAUTHORITY

backtick -E tty { pipeline { tty } tr -d '/dev/tty/' }

foreground { touch $XAUTHORITY }
foreground { xauth remove :$tty }

backtick -E cookie { pipeline { od -An -N16 -tx /dev/urandom } tr -d " " }
foreground { xauth add :$tty MIT-MAGIC-COOKIE-1 $cookie }
Xorg :$tty -keeptty vt$tty -noreset -auth $XAUTHORITY -displayfd 4
