# s6/.s6/sv/dbus-log/run
#!/bin/execlineb -P

importas HOME HOME
envfile ${HOME}/.s6/config/logger.conf
importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES
importas runtime_dir XDG_RUNTIME_DIR
define log_dir ${runtime_dir}/s6/log/dbus

foreground { install -d $log_dir }
s6-log -b -d3 -- ${DIRECTIVES} $log_dir
