# s6/.s6/sv/redshift/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas home HOME
envfile ${home}/.s6/conf/redshift
importas -s args ARGS

redshift ${args}
