# s6/.s6/sv/redshift/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas -S HOME
envfile ${HOME}/.s6/config/redshift.conf
importas -Ssu ARGS ARGS

redshift ${ARGS}
