# s6/.local/share/s6/sv/cdemu-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1

foreground { doas modprobe -a sg sr_mod vhba }

cdemu-daemon
