# s6/.s6/sv/gammastep-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas -S HOME
envfile ${HOME}/.s6/config/night.conf
importas -Ssu ARGS

gammastep ${ARGS}
