# s6/.local/share/s6/sv/mangowc-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
unexport WAYLAND_DISPLAY
importas -S HOME
cd $HOME
/home/senderman/Projects/linux/mango/build/mango -r 4
