# s6/.local/share/s6/sv/foot-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
# Be able to run XWayland apps in terminal
export DISPLAY :0
foot --server
