# s6/.local/share/s6/sv/foot-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1

importas -S XDG_RUNTIME_DIR
elglob MANGO_INSTANCE_SIGNATURE ${XDG_RUNTIME_DIR}/mango-*.sock
export MANGO_INSTANCE_SIGNATURE $MANGO_INSTANCE_SIGNATURE

foot -Ss -lnever -p4
