# s6/.s6/sv/evrouter-srv/finish
#!/usr/bin/execlineb -P

importas DISPLAY DISPLAY
rm /tmp/.evrouter${DISPLAY}
