# s6/.s6/sv/evrouter-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas DISPLAY DISPLAY
foreground { rm /tmp/.evrouter${DISPLAY} }
elglob DEVICES /dev/input/by-id/*
evrouter -f ${DEVICES}
