# s6/.local/share/s6/sv/pipewire-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
export PIPEWIRE_NOTIFICATION_FD 4

pipewire
