# s6/.s6/sv/ranger-dbus/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas HOME HOME

${HOME}/Projects/python/ranger-dbus/run.sh
