# s6/.s6/sv/mangowc-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
unexport WAYLAND_DISPLAY
importas -S HOME
importas -S UID
cd $HOME
mango
