# s6/.local/share/s6/sv/dbus-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas -S DBUS_SESSION_BUS_ADDRESS

dbus-daemon --nofork --session --print-pid 4 --address=${DBUS_SESSION_BUS_ADDRESS} 
