# s6/.local/share/s6/sv/mpd-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas -S XDG_RUNTIME_DIR
foreground { mkdir -p ${XDG_RUNTIME_DIR}/mpd }
mpd --no-daemon --stderr

