# s6/.s6/sv/swayidle/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas -S HOME
envfile ${HOME}/.s6/config/swayidle.conf
importas -Su SCREENSAVER
importas -Su SCREENSAVER_TIMEOUT
foreground { xset s $SCREENSAVER_TIMEOUT }
swayidle -w
timeout $SCREENSAVER_TIMEOUT "${SCREENSAVER}"
lock "${SCREENSAVER}"
