# s6/.local/share/s6/sv/xdg-document-portal-log/run
#!/usr/bin/execlineb -P

s6-logscript
xdg-document-portal
