# s6/.s6/sv/wshowkeys/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas -S HOME
envfile ${HOME}/.s6/config/wshowkeys.conf

multisubstitute {
  importas -Su FONT
  importas -Su SPECIAL
  importas -Su FOREGROUND
  importas -Su BACKGROUND
  importas -Su POS1
  importas -Su POS2
}

wshowkeys -F $FONT -a $POS1 -a $POS2 -s $SPECIAL -f $FOREGROUND -b $BACKGROUND
