# s6/.local/share/s6/sv/battery/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas -S HOME
importas -SD ${HOME}/.local/share XDG_DATA_HOME
envfile ${XDG_DATA_HOME}/s6/config/battery.conf

multisubstitute {
    importas -Su BAT_DEVICE
    importas -Su AC_DEVICE
    importas -Su CRIT_LEVEL
    define power_supply "/sys/class/power_supply"
}

multisubstitute {
    define bat_capacity_file ${power_supply}/${BAT_DEVICE}/capacity
    define ac_online_file ${power_supply}/${AC_DEVICE}/online
}

loopwhilex
    backtick -E bat_lvl { cat $bat_capacity_file }
    backtick -E ac_on { cat $ac_online_file }
    foreground { if { test $ac_on -eq 0 } if { test $bat_lvl -le $CRIT_LEVEL }
        notify-send -u critical "MAYDAY! MAYDAY!" "Low battery! Please, plug in the charger!"
    }
    sleep 5m
