# s6/.local/share/s6/sv/ssh-agent-log/run
#!/usr/bin/execlineb -P

s6-logscript
ssh-agent
