# s6/.local/share/s6/sv/wl-cliphist-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1

importas -S XDG_RUNTIME_DIR

wl-paste --watch cliphist -db-path ${XDG_RUNTIME_DIR}/cliphist store
