# s6/.s6/sv/xss-lock/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas -S HOME
envfile ${HOME}/.s6/config/xss-lock.conf
importas -Su SCREENSAVER

xss-lock $SCREENSAVER
