# s6/.s6/sv/org-freedesktop-filemanager1-common-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1

/usr/lib/file_manager_dbus
