# s6/.local/share/s6/sv/foot-log/run
#!/usr/bin/execlineb -P

s6-logscript
foot
