# s6/.local/share/s6/sv/dbus-log/run
#!/bin/execlineb -P

importas -S HOME
importas -SD ${HOME}/.local/share XDG_DATA_HOME
envfile -I ${XDG_DATA_HOME}/s6/config/logger.conf

multisubstitute {
  define svc dbus # change this to the name of the pipeline
  importas -S XDG_RUNTIME_DIR
  importas -SsCuD "n3 s2000000 T" DIRECTIVES
}

define log_dir ${XDG_RUNTIME_DIR}/s6/log/${svc}

foreground { install -d $log_dir }
s6-log -b -d3 -- ${DIRECTIVES} $log_dir

