# s6/.local/share/s6/sv/gpg-agent-log/run
#!/usr/bin/execlineb -P

s6-logscript
gpg-agent
