# s6/.local/share/s6/sv/pipewire-pulse-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
export PIPEWIRE_PULSE_NOTIFICATION_FD 4

pipewire-pulse
