# s6/.s6/sv/ssh-agent-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas SSH_AUTH_SOCK SSH_AUTH_SOCK

ssh-agent -D -a ${SSH_AUTH_SOCK}

