# s6/.local/share/s6/sv/ssh-agent-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
importas -S HOME
importas -SD ${HOME}/.local/share XDG_DATA_HOME
envfile ${XDG_DATA_HOME}/s6/config/ssh-agent.conf

multisubstitute {
    importas -S SSH_AUTH_SOCK
    importas -Su SSH_AGENT_TIMEOUT
}

ssh-agent -D -t ${SSH_AGENT_TIMEOUT} -a ${SSH_AUTH_SOCK}

