# s6/.s6/sv/wl-cliphist-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
wl-paste --watch cliphist store
