# s6/.local/share/s6/sv/org-freedesktop-filemanager1-common-log/run
#!/usr/bin/execlineb -P

s6-logscript
org-freedesktop-filemanager1-common
